# .config/systemd/user/netns-firewall.target
[Unit]
Description=Active firewall for the user network namespace
Requires=netns-firewall.service
After=netns-firewall.service
