# .config/systemd/user/netns-firewall.target
[Unit]
Description=Ensures firewall for the user network namespace
Requires=netns-firewall.service
After=netns-firewall.service
